# data/init.d/centos/generate_hap2_init_file.sh.in
#!/bin/sh

cat << EOF
#!/bin/sh
# chkconfig: 2345 91 91
# description: HAP2 is a plugin for Hatohol unified manager of
#              monitoring software.

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configurations.
. /etc/sysconfig/network

# Check network is up.
if [ \${NETWORKING} = "no" ]; then
	exit 1
fi

PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin
NAME=@exec_prefix@/libexec/hatohol/hap2/hatohol/hap2_starter.py
exec_prefix=@exec_prefix@
sysconfdir=@expanded_sysconfdir@
CONF=@expanded_sysconfdir@/hatohol/hap2-${CONF_NAME_SUFFIX}.conf
LOG_CONF=@expanded_sysconfdir@/hatohol/hap2-logging.conf
PROG=\$NAME
PIDFILEDIR=/var/run/hatohol
PIDFILE=\${PIDFILEDIR}/hap2_${HAP2_PLUGIN}.pid
OPTIONS="--config \${CONF} --plugin-path \${exec_prefix}/libexec/hatohol/hap2/hatohol/hap2_${HAP2_PLUGIN}.py --pid-file-dir \${PIDFILEDIR} --server-id ${HAP2_PLUGIN} --log-conf \${LOG_CONF}"
PROGNAME=\`basename \$PROG\`
LOCKFILE=/var/lock/subsys/hap2-${CONF_NAME_SUFFIX}

[ -f \$PROG ]|| exit 0

# functions
start() {
	echo -n \$"Starting \$PROGNAME: "
	daemon \$PROG \$OPTIONS
	RETVAL=$?
	echo
	[ \$RETVAL -eq 0 ] && touch \$LOCKFILE
}

stop() {
	echo -n \$"Stopping \$PROGNAME: "
	RETVAL=$?
	kill -SIGTERM \`cat \${PIDFILE}\`
	echo
	[ \$RETVAL -eq 0 ] && rm -f \$LOCKFILE
}

RETVAL=0
case "\$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	status)
		status \$NAME
		;;
	*)
		echo \$"Usage: \$PROGNAME {start|stop|restart|status}" >&2
		exit 1
		;;
esac
exit \$RETVAL
EOF
